# systemd/super-arr-stack-region-rotate.service
[Unit]
Description=Rotate Super Arr Stack VPN region
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
Environment=APP_ROOT=/opt/super-arr-stack
ExecStart=/opt/super-arr-stack/scripts/rotate-vpn-region.sh
